# deploy/systemd/pipeworks-namegen-api.service.example
[Unit]
Description=Pipeworks Namegen API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pipeworks
Group=pipeworks
WorkingDirectory=/opt/pipeworks-namegen-api/current
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/pipeworks-namegen-api/venv/bin/python -m pipeworks_namegen_api.webapp.api --config /etc/pipeworks/namegen-api/server.ini
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/pipeworks-namegen-api /var/log/pipeworks-namegen-api

[Install]
WantedBy=multi-user.target
